# debian/postinst
#! /bin/sh
# postinst script for cpufreqd
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#
# quoting from the policy:
#     Any necessary prompting should almost always be confined to the
#     post-installation script, and should be protected with a conditional
#     so that unnecessary prompting doesn't happen if a package's
#     installation fails and the `postinst' is called with `abort-upgrade',
#     `abort-remove' or `abort-deconfigure'.

# Source debconf library.
. /usr/share/debconf/confmodule

case "$1" in
    configure)
        if [ -d /sys/devices/system/cpu/cpu0/cpufreq ] || [ -f /proc/cpufreq ] ; then
            
          [ -z ${CPUFREQD_PM} ] || unset CPUFREQD_PM

          if [ -d /proc/pmu ]; then
            CPUFREQD_PM=pmu
          elif [ -f /proc/apm ]; then
            CPUFREQD_PM=apm
          elif [ -d /proc/acpi ]; then
            CPUFREQD_PM=acpi
          fi
          

          if [ ${CPUFREQD_PM} ]; then
            # setup a very simple default configuration file
            sed -i -e "s/^\(pm_type=\)acpi\(.*\)/\1$CPUFREQD_PM\2/" /etc/cpufreqd.conf

            update-rc.d cpufreqd defaults > /dev/null

            if [ -x /usr/sbin/invoke-rc.d ] ; then
                  invoke-rc.d cpufreqd start
            else
                /etc/init.d/cpufreqd start
            fi

          fi

        fi

        # warn the user about not finding cpufreq kernel interface

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

db_stop

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

